# debian/loctrkd.wsgateway.service
[Unit]
Description=GPS303 Websocket Gateway Service
PartOf=loctrkd.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/loctrkd
ExecStart=python3 -m loctrkd.wsgateway $OPTIONS
KillSignal=INT
Restart=on-failure
StandardOutput=journal
StandardError=inherit
User=loctrkd
Group=loctrkd

[Install]
WantedBy=loctrkd.target
